# config/systemd/user/yay-cache-clean.service
[Unit]
Description=Clean yay cache

[Service]
Type=oneshot
ExecStart=/usr/bin/find %h/.cache/yay -mindepth 1 -maxdepth 1 -mtime +30 -exec rm -rf {} +
